# packaging/linux/dotvault.service
[Unit]
Description=dotvault — Vault-to-file secret synchronisation daemon
Documentation=https://github.com/goodtune/dotvault

[Service]
Type=notify
NotifyAccess=main
TimeoutStartSec=300
WatchdogSec=120

ExecStart=/usr/bin/dotvault run

# SIGHUP is the daemon's reload trigger: it re-reads ~/.dotvault-token and
# immediately re-runs the config loader, applying the dynamic sections
# (rules, enrolments, sync interval, remote_config) in place. Static
# sections (vault, web, agent, observability, bypass_system_config) are
# logged as needing a restart.
ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal

EnvironmentFile=-%h/.config/dotvault/env
EnvironmentFile=-%h/.config/dotvault.env
EnvironmentFile=-/etc/default/dotvault
EnvironmentFile=-/etc/sysconfig/dotvault

# Only sandbox directives that work in a non-root user manager. The
# Protect*/Private* knobs need CAP_SETPCAP / CAP_SYS_ADMIN and would
# fail the spawn at step CAPABILITIES with exit 218.
NoNewPrivileges=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=default.target
