# packaging/linux/dotvault-token-watch.path
[Unit]
Description=Watch ~/.vault-token for changes
PartOf=dotvault.service

[Path]
PathChanged=%h/.vault-token
Unit=dotvault-token-watch.service

[Install]
WantedBy=default.target
